# supervisord.conf
[supervisord]
nodaemon=true

[program:redis]
command=/usr/bin/redis-server

[program:node]
command=node index.js
directory=/app
autostart=true
autorestart=true
